# playtable.conf
[program:playtable]
command=/usr/bin/wine %(ENV_VPT_DIR)s/Playtable.exe -server
priority=100
stdout_logfile=/var/log/playtable.stdout.log
stderr_logfile=/var/log/playtable.stderr.log